# apps/hello-simple/hello.conf
[program:hello]
command=java -jar /var/hello/hello-world-1.0.0-SNAPSHOT-runner.jar
autostart=true
autorestart=true
stderr_logfile=/var/log/hello.err.log
stdout_logfile=/var/log/hello.out.log
